# linux_pkg/deb/deb_files/DEBIAN/prerm
# !/bin/sh
systemctl stop svas
systemctl disable svas
