# linux_pkg/deb/dpkg/DEBIAN/prerm
# !/bin/sh
systemctl stop svas
systemctl disable svas
rm -f /etc/systemd/system/multi-user.target.wants/svas.service